# infrastructure/systemd/fm-sync-stages-test.service
# Stages test server (port 3002) – lookup user by Blackbaud ID, see checklist + stage/substage.
# Optional: only enable if you need this UI in production.
# Copy to /etc/systemd/system/ and edit paths if your app is not in /opt/fm-sync.
#
#   sudo cp infrastructure/systemd/fm-sync-stages-test.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now fm-sync-stages-test

[Unit]
Description=FM Sync stages test server (port 3002)
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/fm-sync
# So systemd finds node (tsx uses env node). Add your node bin path if you use nvm: e.g. /home/YOU/.nvm/versions/node/v20.x.x/bin
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=STAGES_TEST_PORT=3002
# Absolute path required by systemd. If you symlink /opt/fm-sync -> your repo, this works.
ExecStart=/opt/fm-sync/node_modules/.bin/tsx packages/processor/src/stages-test-server.ts
Restart=on-failure
RestartSec=30

# Run as unprivileged user if desired
# User=fm-sync
# Group=fm-sync

[Install]
WantedBy=multi-user.target
